# deploy/breakyourstudy.service
[Unit]
Description=BreakYourStudy Gradio app
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/breakYourStudy
# 固定本机监听端口（与 deploy/breakyourstudy.conf 中 upstream 一致）。load_dotenv(override=False) 时不会被子级 .env 覆盖。
Environment=PORT=17860
Environment=GRADIO_SERVER_NAME=127.0.0.1
# app.py uses load_dotenv(); .env must live in WorkingDirectory
ExecStart=/root/breakYourStudy/.venv/bin/python /root/breakYourStudy/app.py
Restart=on-failure
RestartSec=5
# Gradio can be slow to stop on SIGTERM during active sessions
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
